# emos-boot-manager.service
[Unit]
Description=EMOS Boot Mode Manager
...